# hive/share/systemd/restart-asic-cron.service
[Unit]
Description=restart-asic-cron

[Service]
Type=simple
ExecStart=/sbin/reboot

[Install]
WantedBy=multi-user.target
